# server/email-service.service

[Unit]
Description=FreePBX CRM Email Service
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/html/crm/server
ExecStart=/usr/bin/node email-service.js
Restart=always
RestartSec=10
Environment=NODE_ENV=production
Environment=EMAIL_PORT=3002

[Install]
WantedBy=multi-user.target
